# docker/culvertsense-secrets.service
[Unit]
Description=Load CulvertSense secrets from Azure Key Vault
After=network-online.target
Wants=network-online.target
Before=docker.service culvertsense.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/load-secrets.sh

[Install]
WantedBy=multi-user.target
